# Linux22/otus-linux-network/provision/knock.service
[Unit]
Description=knocking port

[Service]
Type=simple
ExecStart=/etc/iptables/knock.sh

[Install]
WantedBy=multi-user.target
