# Linux26/provision/roles/db/files/backup.service
[Unit]
Description=Backup of etc

[Service]
Type=simple
ExecStart=/vagrant/provision/backup.sh

[Install]
WantedBy=multi-user.target
